# etc/xinetd.d/telnets
# default: on
# description: TelnetS is a secure SSL login server.
service telnets
{
        socket_type     = stream
        protocol        = tcp
        wait            = no
        user            = root
        flags           = NAMEINARGS
        server          = /usr/sbin/stunnel
        server_args     = . /etc/in.telnetd
        disable         = no
}
